# com.cpcc.kiosk-ce-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cpcc.kiosk-ce-refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/frazier/Documents/Administrative/AdvisingAndCareerDay/degree-sheets/kiosk/refresh.sh</string>
  </array>
  <!-- Run weekly: Mondays at 6:00 AM -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>6</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/frazier/Documents/Administrative/AdvisingAndCareerDay/degree-sheets/ce-refresh.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/frazier/Documents/Administrative/AdvisingAndCareerDay/degree-sheets/ce-refresh.err.log</string>
</dict>
</plist>
